# systemd/freenect2cam.service
[Unit]
Description=Freenect2cam

[Service]
ExecStart=/usr/bin/freenect2cam

[Install]
WantedBy=multi-user.target